# supervisor.conf
[program:homepage]
command=docker-compose -f production.yml up
directory=/home/homepage/site
redirect_stderr=true
autostart=true
autorestart=true
priority=10
